# run/portable-service/tclip.service.in
[Unit]
Description=The tclip service
Documentation=https://github.com/tailscale-dev/tclip

[Service]
DynamicUser=yes
StateDirectory=tclip
Environment=DATA_DIR=/var/lib/private/tclip/data
ExecStart=@web@/bin/web
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
