# run/portable-service/tclip.service.in
[Unit]
Description=The Infinipaste service
Documentation=https://github.com/tailscale-dev/infinipaste

[Service]
DynamicUser=yes
StateDirectory=infinipaste
Environment=DATA_DIR=/var/lib/private/infinipaste/data
ExecStart=@web@/bin/web
Restart=always
RestartSec=30s

[Install]
WantedBy=multi-user.target
